# deploy/kafo.service.example
[Unit]
Description=Kafo ASP.NET Core Web Application
After=network.target clamav-daemon.service
Wants=clamav-daemon.service

[Service]
Type=notify
User=kafo
Group=kafo
WorkingDirectory=/var/www/kafo/current
EnvironmentFile=/etc/kafo/kafo.env
ExecStart=/usr/bin/dotnet /var/www/kafo/current/Kafo.Web.dll
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
SyslogIdentifier=kafo

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kafo /var/log/kafo
UMask=0027

[Install]
WantedBy=multi-user.target
